# scripts/com.ulos.telegram-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ulos.telegram-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/grzegorzgoldyn/dev/ul-os-telegram-bot/.venv/bin/python</string>
        <string>-m</string>
        <string>app.main</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/grzegorzgoldyn/dev/ul-os-telegram-bot</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/grzegorzgoldyn/dev/ul-os-telegram-bot/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/grzegorzgoldyn/dev/ul-os-telegram-bot/logs/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/grzegorzgoldyn/dev/ul-os-telegram-bot/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
